# src/fs/etc/services.d/mount/run
#!/usr/bin/with-contenv sh
. "/usr/bin/variables"

umask 022

mount_command="mount"

exec $mount_command